# srcs/nginx/supervisor_conf.ini
[supervisord]
nodaemon=true
[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
[program:influxdb]
command=/influxdb2-2.0.4-linux-amd64/influxd
autostart=true
autorestart=true